# deploy/tender-agents.service
[Unit]
Description=FeedBackTalk Tender Leads dashboard
After=network.target

[Service]
Type=simple
User=esasim
Group=esasim
WorkingDirectory=/home/esasim/agents
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/esasim/agents/.env
ExecStart=/home/esasim/agents/.venv/bin/tender-leads serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
